# cron.d/flexget
# m h dom mon dow user  command
0 *   * * *   abc    /usr/local/bin/flexget --loglevel debug -c /config/config.yml execute --now 2>&1 | logger -t FlexGet
